# runit/tty-no-cursor-blink/run
#!/bin/sh
exec echo 0 > /sys/class/graphics/fbcon/cursor_blink
